# ec20-feishu-cleanup.timer
[Unit]
Description=Daily EC20 Feishu message retention cleanup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
